# roles/valheim/files/valheim-update.timer
[Unit]
Description=Restart the server every night

[Timer]
OnCalendar=*-*-* 00:00:00
Unit=valheim-update.service

[Install]
WantedBy=timers.target
