# roles/valheim/files/valheim-update.service
[Unit]
Description=Valheim Update

[Service]
Type=oneshot
ExecStart=/home/ubuntu/do-update.sh
User=ubuntu
WorkingDirectory=/home/ubuntu
