# roles/valheim-backup/files/valheim-backup.timer
[Unit]
Description=Schedule a backup every 1 hour

[Timer]
OnCalendar=*-*-* *:00:00
Unit=valheim-backup.service

[Install]
WantedBy=timers.target
