# roles/valheim-backup/files/valheim-backup.service
[Unit]
Description=Valheim Backup

[Service]
Type=oneshot
ExecStart=/home/ubuntu/make-backup.sh
User=ubuntu
WorkingDirectory=/home/ubuntu
